# supervisord.conf
; Supervisord 配置文件
; 用于管理容器内的多个进程

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisor/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

; ==================== 程序配置 ====================

; Nginx Web 服务器
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
process_name=%(program_name)s
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=false
stdout_logfile=/var/log/nginx/nginx_access.log
stderr_logfile=/var/log/nginx/nginx_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=

; Cron 定时任务
[program:cron]
command=/usr/sbin/cron -f
process_name=%(program_name)s
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=false
stdout_logfile=/var/log/cron/cron_stdout.log
stderr_logfile=/var/log/cron/cron_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=

; Stock Monitor Web 应用
[program:stock_monitor]
command=/usr/bin/python3 start_app.py
process_name=%(program_name)s
numprocs=1
directory=/root/stock_monitor
autostart=true
autorestart=true
startsecs=5
startretries=2
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/stock_monitor/web_stdout.log
stderr_logfile=/var/log/stock_monitor/web_stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1",TZ="Asia/Shanghai",PYTHONPATH="/root/crawlers/spiders:/root"

; SSH 服务
[program:sshd]
command=/usr/sbin/sshd -D
process_name=%(program_name)s
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=false
stdout_logfile=/var/log/ssh/sshd_stdout.log
stderr_logfile=/var/log/ssh/sshd_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=

; ==================== 事件监听器（可选） ====================

; [eventlistener:crashmail]
; command=/usr/bin/python3 /root/scripts/crash_report.py
; events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
; stderr_logfile=/var/log/supervisor/crashmail_stderr.log
; stdout_logfile=/var/log/supervisor/crashmail_stdout.log

; ==================== 组配置（可选） ====================

[group:stock_monitor]
programs=stock_monitor
priority=100

[group:system]
programs=nginx,cron,sshd
priority=50
